# setup/scripts/doorphoneserver.service
# Servizio systemd per DoorPhoneServer.
# ProtectSystem=false: necessario per accesso a device audio, GPIO e /var/log.

[Unit]
Description=DoorPhoneServer Service
Requires=network.target sound.target
After=network-online.target sound.target
Wants=network-online.target

[Service]
User=doorphoneserver
Group=doorphoneserver
Type=simple
WorkingDirectory=/home/doorphoneserver
ExecStartPre=/bin/sleep 15
ExecStart=/home/doorphoneserver/bin/doorphoneserver -config=/home/doorphoneserver/doorphoneserver.xml
ExecStop=/bin/kill -s SIGTERM $MAINPID
Restart=always
RestartSec=10

# PrivateTmp=true: processo isolato in /tmp privata
# ProtectSystem=false: necessario per accesso in scrittura a /var/log e device di sistema
PrivateTmp=true
ProtectSystem=false

[Install]
WantedBy=multi-user.target
